# etc/udev/rules.d/80-keychron.rules
SUBSYSTEMS=="input", ATTRS{name}=="Keychron K2", RUN+="echo 2 | tee /sys/module/hid_apple/parameters/fnmode"
